# mymcp.service
# systemd --user unit for mymcp (no sudo).
# Install:   task install-unit   (copies this to ~/.config/systemd/user/)
# Configure: edit ~/.config/mymcp/mymcp.env  (WORKSPACE, ADDR)
# Start:     systemctl --user start mymcp
[Unit]
Description=mymcp — MCP tool server (mandatory mTLS, per-CN capabilities)
Documentation=https://github.com/brywil/mymcp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/mymcp/mymcp.env
ExecStart=%h/.local/bin/mymcp serve --config-dir %h/.config/mymcp --addr ${ADDR} --workspace ${WORKSPACE}
Restart=on-failure
RestartSec=2

# Modest hardening only: mymcp's run_command is remote code execution by design,
# so heavy filesystem sandboxing would fight its purpose. The real protection is
# mandatory mTLS + deny-by-default per-CN capabilities. We still keep the system
# read-only and give it a private /tmp.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=%h
PrivateTmp=true

[Install]
WantedBy=default.target
